# config/etc/systemd/system/docker.service.d/25-bridge.conf
[Service]
ExecStartPost=sh -c 'iptables -C DOCKER-USER -i br0 -o br0 -j ACCEPT >/dev/null 2>&1 || iptables -I DOCKER-USER -i br0 -o br0 -j ACCEPT'
