# config/etc/local.d/zram.stop
#!/usr/bin/env bash

swapoff /dev/zram0

zramctl --reset /dev/zram0

modprobe --remove zram
